# zxmaster.service
[Unit]
Description=Zxmaster-Bot

[Service]
ExecStart=/home/pi/Documentos/Zxmaster-Bot/zxmaster_script.sh
User=pi
Group=pi

[Install]
WantedBy=network-online.target
